# com.moviehub.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  MovieHub 이벤트 자동 갱신 스케줄러 (launchd)

  설치:
    cp com.moviehub.refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.moviehub.refresh.plist

  즉시 실행 테스트:
    launchctl start com.moviehub.refresh

  제거:
    launchctl unload ~/Library/LaunchAgents/com.moviehub.refresh.plist
    rm ~/Library/LaunchAgents/com.moviehub.refresh.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moviehub.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/sysilion/proj/moviehub/tools/refresh_events.py</string>
    </array>

    <!-- 매 정각 실행 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/sysilion/proj/moviehub</string>

    <key>StandardOutPath</key>
    <string>/Users/sysilion/proj/moviehub/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sysilion/proj/moviehub/server.log</string>

    <!-- Mac이 절전 모드에서 깨어날 때 지난 실행을 즉시 보충 실행 -->
    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
